# arm32/ubuntu-rootfs/ubuntu_18.04/lib/systemd/system/advantech.service
[Unit]
Description=Advantech's program running
After=network.target auditd.service

[Service]
ExecStart=/etc/init.d/advantech.sh
Type=notify

[Install]
WantedBy=multi-user.target
