# arm32/ubuntu-rootfs/ubuntu_16.04/lib/systemd/system/advantech.service
[Unit]
Description=Just for advantech's program running

[Service]
ExecStart=/etc/init.d/advantech.sh

Type=forking

[Install]
WantedBy=multi-user.target
